# andon-radio.service
[Unit]
Description=Andon FM headless radio player
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=aktasch
ExecStart=/usr/bin/python3 /home/aktasch/andon-radio/radio.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
